# containers/stoat/stoat-redis.container
[Container]
Image=docker.io/eqalpha/keydb
Pod=stoat.pod

[Service]
Restart=always
