# containers/stoat/stoat-web.container
[Container]
EnvironmentFile=.env.web
Image=ghcr.io/revoltchat/client:master
Pod=stoat.pod

[Service]
Restart=always
