# containers/stoat/stoat-january.container
[Container]
Image=ghcr.io/revoltchat/january:20250930-2
Pod=stoat.pod
Volume=./Revolt.toml:/Revolt.toml

[Service]
Restart=always
